# wgman.target
[Unit]
Description=WireGuard Tunnels via wgman
